# debug_broken_app_stack/environment/supervisord.conf
[program:redis]
command=redis-server /etc/redis/redis.conf
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/redis.log
stderr_logfile=/var/log/supervisor/redis_err.log

[program:init_db]
command=python3 /app/init_db.py
autostart=true
autorestart=false
startsecs=0
priority=15
stdout_logfile=/var/log/supervisor/init_db.log
stderr_logfile=/var/log/supervisor/init_db_err.log

[program:gunicorn]
command=gunicorn -c /app/gunicorn_config.py app:app
directory=/app
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/supervisor/gunicorn.log
stderr_logfile=/var/log/supervisor/gunicorn_err.log

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=30
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx_err.log
